# utils/supervisor/dede.conf
[program:dede]

command=dfm_path/deepdetect/build/main/dede

user=dfm

autostart=true

autorestart=true

directory=dfm_path/deepdetect
