# utils/supervisor/dfm.conf
[program:dfm]

command=/usr/bin/python dfm_path/dfm/server.py

user=dfm

autostart=true

autorestart=true

redirect_stderr=true

directory=dfm_path/dfm
