# utils/supervisor/es.conf
[program:es]

command=dfm_path/elasticsearch/bin/elasticsearch

user=dfm

autostart=true

autorestart=true

directory=dfm_path/elasticsearch
