# utils/supervisor/kb.conf
[program:kibana]

command=dfm_path/kibana/bin/kibana

user=dfm

autostart=true

autorestart=true

directory=dfm_path/kibana
